# example/transplied/db.service
[Unit]
After=network.target
Description=Database

[Service]
ExecStart=/usr/bin/db-server

[Install]
WantedBy=multi-user.target
